# ops/sv-uwsgi.conf
[program:uwsgi-gsevolapp]
command=/root/Env/gsevol/bin/uwsgi /root/gsevol-web/ops/uwsgi.ini
autostart=true
autorestart=true
stdout_logfile=/root/gsevol-web/uwsgi.log
redirect_stderr=true
exitcodes=0
